# supervisord-user.conf
[program:sshd]
priority=30
command=/usr/sbin/sshd -D


